# tdx/mkosi.extra/etc/systemd/system/tdxs.socket
[Unit]
Description=TDXS quote socket

[Socket]
ListenStream=/var/tdxs.sock
SocketMode=0660
SocketUser=root
SocketGroup=tdx
Accept=false

[Install]
WantedBy=sockets.target
